# packages/mediacenter/plexht/system.d/plexht-autostart.service
[Unit]
Description=OpenPHT user autostart script
Before=plexht.service
After=network-online.target graphical.target
ConditionPathExists=/storage/.config/autostart.sh

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=-/bin/sh -c ". /etc/profile; exec /bin/sh /storage/.config/autostart.sh"
RemainAfterExit=yes

[Install]
WantedBy=plexht.service
