# packages/mediacenter/plexht/system.d/storage-.kodi.mount
[Unit]
DefaultDependencies=no
Conflicts=umount.target
After=systemd-tmpfiles-setup.service
Before=sysinit.target umount.target

[Mount]
What=/storage/.plexht
Where=/storage/.kodi
Type=none
Options=bind

[Install]
WantedBy=sysinit.target
